# mqtt-alive-daemon.service
[Unit]
Description=MQTT Alive Daemon
After=network.target

[Service]
ExecStart=/path/to/your/mqtt-alive-daemon
Restart=always
User=yourusername
Environment="MQTT_PASSWORD=your_mqtt_password"

[Install]
WantedBy=multi-user.target
